# prebuilt/common/etc/cron/cron.daily/00drop_caches
#!/system/bin/sh

. /system/etc/cron/cron.conf

if test $DAILY_DROP_CACHES -eq 1; then
	sync
	echo 3 > /proc/sys/vm/drop_caches
fi