# prebuilt/common/etc/cron/cron.weekly/01clear_cache
#!/system/bin/sh

. /system/etc/cron/cron.conf

if test $WEEKLY_CLEAR_CACHE -eq 1; then
	cache -rmdata
fi